# mqtt-sensors/bme280.conf.example
[program:bme280]
command=sudo /home/pi/RPiSensorBox/mqtt-sensors/env/bin/python /home/pi/RPiSensorBox/mqtt-sensors/bme280.py
